# roles/system/files/etc/profile.d/electron-wayland.sh
# Use wayland for electron, see
# https://github.com/electron/electron/blob/main/docs/api/environment-variables.md#electron_ozone_platform_hint-linux
export ELECTRON_OZONE_PLATFORM_HINT="auto"
