# roles/system/files/etc/profile.d/xdg-data-dirs.sh
# to make junction work
export XDG_DATA_DIRS="/usr/local/share/:/usr/share/:${XDG_DATA_DIRS}"
