# roles/plasma/files/etc/profile.d/gtk-file-dialog-portal.sh
# consistent file dialogs for gtk apps under plasma
export GTK_USE_PORTAL=1
export GDK_DEBUG=portals
